# client/einar.service
[Unit]
After=network.target

[Service]
ExecStart=LOKIHOME/MYNAME-onboot.sh

[Install]
WantedBy=default.target